# systemd/pve-b2-age-prune.timer
[Unit]
Description=Daily prune for PVE B2 Age Backups
Documentation=https://github.com/LukasStrickler/proxmox-b2-offsite-backup

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
