# systemd/pve-b2-age-hostconfig.service
[Unit]
Description=PVE B2 Age Host Config Backup
Documentation=https://github.com/LukasStrickler/proxmox-b2-offsite-backup

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pve-b2-age-hostconfig.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pve-b2-age-hostconfig

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/tmp /run/lock
ReadOnlyPaths=/etc /var/lib/pve-cluster /root /etc/pve-b2-age-backup
